# deploy/freqtrade.service
[Unit]
Description=Freqtrade NostalgiaForInfinity (dry-run)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/botbot2
EnvironmentFile=/home/ubuntu/botbot2/.env
ExecStart=/home/ubuntu/botbot2/.venv/bin/freqtrade trade --userdir /home/ubuntu/botbot2/user_data --config /home/ubuntu/botbot2/user_data/config_nfi.json --strategy NostalgiaForInfinityX7
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Memory safety on a 2 GB VPS (relies on swap; OOM-restart instead of freezing host)
MemoryMax=1700M
MemorySwapMax=2G

[Install]
WantedBy=multi-user.target
